# RoboterSteuerung/Robotersteuerung_v1_inc_AS/Verwaltungsschale/install/control.service
[Unit]
Description=i40AS Control

[Service]
User=USR
Group=www-data
WorkingDirectory=DIR
Environment="PATH=DIR/projectenv/bin"
ExecStart=DIR/projectenv/bin/python3.5 DIR/control_start.py

[Install]
WantedBy=multi-user.target
